# service-auditor.cron
* * * * * root /usr/local/bin/php /var/www/html/public/run-check.php >> /var/log/service-auditor-cron.log 2>&1
